# rpi_thermo_chick/rpi-thermo-chick.service
[Unit]
Description=rpi-thermo-chick-service
After=network.target

[Service]
ExecStart=rpi-thermo-chick

[Install]
WantedBy=multi-user.target